# systemd/botifler.timer
[Unit]
Description=Run Botifler Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
